# scripts/strike2.conf
[program:strike2]
command=/bin/bash -c "java -jar /srv/strike2/strike2-0.3.0-standalone.jar"
autostart=true
autorestart=true